# 60-ipu6-tgl-adl.rules
KERNEL=="video[0-9]*", RUN+="ipu-setlink %p"

